# desktop/macos/weisile-link.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>cn.vsle.weisile-link</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/WeisileLink.app/Contents/MacOS/WeisileLink</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WEISILE_LINK_HOST</key>
    <string>127.0.0.1</string>
    <key>WEISILE_LINK_PORT</key>
    <string>20111</string>
    <key>TRAINER_WS_PORT</key>
    <string>8766</string>
    <key>WEISILE_TRANSPORT</key>
    <string>wifi</string>
    <key>LOG_LEVEL</key>
    <string>INFO</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>~/Library/Logs/WeisileLink/weisile-link.out.log</string>
  <key>StandardErrorPath</key>
  <string>~/Library/Logs/WeisileLink/weisile-link.err.log</string>
</dict>
</plist>
